# config/global/etc/systemd/system/cni-restart.service
[Unit]
Description=Reinitialise CNI bridge after restart
Before=kubelet.service strongswan.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/cni-restart.sh
StandardOutput=syslog
StandardError=syslog

[Install]
WantedBy=basic.target